# fbcpCropped.service
[Unit]
Description=juj fbcpCropped driver for Freeplay CM3
After=basic.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=2
User=root
ExecStart=/usr/local/bin/fbcpCropped

[Install]
WantedBy=multi-user.target
