# fbcpFilled.service
[Unit]
Description=juj fbcpFilled driver for Freeplay CM3
After=basic.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=2
User=root
ExecStart=/usr/local/bin/fbcpFilled

[Install]
WantedBy=multi-user.target
